# stack/systemd/llama-jina.service
[Unit]
Description=llama.cpp Jina embedding server
Documentation=file:/home/lexa/an/git/local_llm_setup/stack/systemd/README.md
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/lexa/an/git/local_llm_setup
ExecStart=/home/lexa/an/git/local_llm_setup/stack/llama_cpp/start_jina.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
